# deploy/mesh-worker@.service
[Unit]
Description=AI Mesh Worker (%i)
After=network-online.target wg-quick@wg0.service
Wants=network-online.target
Requires=wg-quick@wg0.service

[Service]
Type=simple
User=mesh-worker
Group=mesh
SupplementaryGroups=sam
WorkingDirectory=/opt/mesh-worker
EnvironmentFile=/etc/mesh-worker/%i.env
EnvironmentFile=-/etc/mesh-worker/common.env
ExecStart=/opt/mesh-worker/venv/bin/python -m src.router.worker_client
Restart=always
RestartSec=10
UMask=0002

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/mesh-worker/.mesh
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mesh-worker-%i

[Install]
WantedBy=multi-user.target
